# root/etc/services.d/authserver/run
#!/usr/bin/with-contenv bash

cd /app/server/bin || exit

exec \
	s6-setuidgid abc ./authserver
